# udev/50-udev.rules
# do not edit this file, it will be overwritten on update

SUBSYSTEM=="block", SYMLINK{unique}+="block/%M:%m"
SUBSYSTEM!="block", SYMLINK{unique}+="char/%M:%m"

KERNEL=="pty[pqrstuvwxyzabcdef][0123456789abcdef]", GROUP="tty", MODE="0660"
KERNEL=="tty[pqrstuvwxyzabcdef][0123456789abcdef]", GROUP="tty", MODE="0660"
KERNEL=="ptmx",     GROUP="tty", MODE="0666"
KERNEL=="tty",      GROUP="tty", MODE="0666"
KERNEL=="tty[0-9]*",    GROUP="tty", MODE="0620"
KERNEL=="console",    MODE="0600"
KERNEL=="vcs|vcs[0-9]*|vcsa|vcsa[0-9]*", GROUP="tty"

# serial
KERNEL=="tty[A-Z]*[0-9]|pppox[0-9]*|ircomm[0-9]*|noz[0-9]*|rfcomm[0-9]*", GROUP="dialout"
KERNEL=="ppp",      MODE="0600"
KERNEL=="mwave",    GROUP="dialout"
KERNEL=="hvc*|hvsi*",   GROUP="dialout"

# virtio serial / console ports
KERNEL=="vport*", ATTR{name}=="?*", SYMLINK+="virtio-ports/$attr{name}"

# mem
KERNEL=="null|zero|full|random|urandom", MODE="0666"
KERNEL=="mem|kmem|port|nvram",  GROUP="kmem", MODE="0640"

# input
KERNEL=="mouse*|mice|event*", MODE="0640"
KERNEL=="ts[0-9]*|uinput",  MODE="0640"
KERNEL=="js[0-9]*",   MODE="0644"

# video4linux
SUBSYSTEM=="video4linux", GROUP="video"
KERNEL=="vttuner*",   GROUP="video"
KERNEL=="vtx*|vbi*",    GROUP="video"
KERNEL=="winradio*",    GROUP="video"

# graphics
KERNEL=="agpgart",    MODE="0600", GROUP="video"
KERNEL=="pmu",      GROUP="video"
KERNEL=="nvidia*|nvidiactl*", GROUP="video"
SUBSYSTEM=="graphics",    GROUP="video"
SUBSYSTEM=="drm",   GROUP="video"

# sound
SUBSYSTEM=="sound",   GROUP="audio"
KERNEL=="mixer0",   SYMLINK+="mixer"

# DVB (video)
SUBSYSTEM=="dvb", GROUP="video"

# FireWire (firewire-core driver: IIDC devices, AV/C devices)
SUBSYSTEM=="firewire", ATTR{units}=="*0x00a02d:0x00010*", GROUP="video"
SUBSYSTEM=="firewire", ATTR{units}=="*0x00a02d:0x010001*", GROUP="video"

# FireWire (deprecated dv1394 and video1394 drivers)
KERNEL=="dv1394-[0-9]*",  NAME="dv1394/%n", GROUP="video"
KERNEL=="video1394-[0-9]*", NAME="video1394/%n", GROUP="video"

# libusb device nodes
SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", MODE="0664"

# printer
KERNEL=="parport[0-9]*",  GROUP="lp"
SUBSYSTEM=="printer",   KERNEL=="lp*", GROUP="lp"
SUBSYSTEM=="ppdev",   GROUP="lp"
SUBSYSTEM=="usb",   KERNEL=="lp*", SYMLINK+="usb%k", GROUP="lp"
KERNEL=="lp[0-9]*",   GROUP="lp"
KERNEL=="irlpt[0-9]*",    GROUP="lp"
# hplip and cups 1.4+ use raw USB devices, so permissions should be similar to
# the ones from the old usblp kernel module
SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", ENV{ID_USB_INTERFACES}=="", IMPORT{program}="usb_id --export %p"
SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", ENV{ID_USB_INTERFACES}==":0701*:", GROUP="lp", MODE="0664"

# block
SUBSYSTEM=="block", GROUP="disk"

# floppy
KERNEL=="fd[0-9]", GROUP="floppy"
KERNEL=="fd[0-9]", ACTION=="add", ATTRS{cmos}=="?*", RUN+="create_floppy_devices -c -t $attr{cmos} -m %M -M 0640 -G floppy $root/%k"
KERNEL=="hd*", SUBSYSTEMS=="ide", ATTRS{media}=="floppy", OPTIONS+="all_partitions"

# cdrom
SUBSYSTEM=="block", KERNEL=="sr[0-9]*", SYMLINK+="scd%n", GROUP="cdrom"
SUBSYSTEM=="block", KERNEL=="hd*", SUBSYSTEMS=="ide", ATTRS{media}=="cdrom", GROUP="cdrom"
SUBSYSTEM=="scsi_generic", SUBSYSTEMS=="scsi", ATTRS{type}=="4|5", GROUP="cdrom"
KERNEL=="pktcdvd[0-9]*", GROUP="cdrom"
KERNEL=="pktcdvd", GROUP="cdrom"

# tape
KERNEL=="ht[0-9]*|nht[0-9]*", GROUP="tape"
KERNEL=="pt[0-9]*|npt[0-9]*|pht[0-9]*", GROUP="tape"
SUBSYSTEM=="scsi_generic|scsi_tape", SUBSYSTEMS=="scsi", ATTRS{type}=="1|8", GROUP="tape"

# block-related
KERNEL=="sch[0-9]*", GROUP="disk"
SUBSYSTEM=="scsi_generic", SUBSYSTEMS=="scsi", ATTRS{type}=="0", GROUP="disk"
KERNEL=="pg[0-9]*", GROUP="disk"
KERNEL=="qft[0-9]*|nqft[0-9]*|zqft[0-9]*|nzqft[0-9]*|rawqft[0-9]*|nrawqft[0-9]*", GROUP="disk"
KERNEL=="rawctl", GROUP="disk"
SUBSYSTEM=="raw", KERNEL=="raw[0-9]*", GROUP="disk"
SUBSYSTEM=="aoe", GROUP="disk", MODE="0220"
SUBSYSTEM=="aoe", KERNEL=="err", MODE="0440"

# network
KERNEL=="tun",      MODE="0666"
KERNEL=="rfkill",   MODE="0644"

# CPU
KERNEL=="cpu[0-9]*",    MODE="0444"

# miscellaneous
KERNEL=="fuse",     MODE="0666"
SUBSYSTEM=="rtc", DRIVERS=="rtc_cmos", SYMLINK+="rtc"
KERNEL=="mmtimer",    MODE="0644"
KERNEL=="rflash[0-9]*",   MODE="0400"
KERNEL=="rrom[0-9]*",   MODE="0400"

# do not delete static device nodes
ACTION=="remove", NAME=="", TEST=="/lib/udev/devices/%k", OPTIONS+="ignore_remove"
ACTION=="remove", NAME=="?*", TEST=="/lib/udev/devices/$name", OPTIONS+="ignore_remove"
#Mapping for USB Power Adaptor 1
 
#Mapping for USB Power Adaptor 1
BUS=="usb", ID=="9-1.2.1:1.0", SYMLINK+="gsmmd1"
BUS=="usb", ID=="9-1.2.2:1.0", SYMLINK+="gsmmd2"
BUS=="usb", ID=="9-1.2.3:1.0", SYMLINK+="gsmmd3"
BUS=="usb", ID=="9-1.2.4:1.0", SYMLINK+="gsmmd4"
BUS=="usb", ID=="9-1.2.5:1.0", SYMLINK+="gsmmd5"
BUS=="usb", ID=="9-1.2.6:1.0", SYMLINK+="gsmmd6"
BUS=="usb", ID=="9-1.2.7.1:1.0", SYMLINK+="gsmmd7"
BUS=="usb", ID=="9-1.2.7.2:1.0", SYMLINK+="gsmmd8"
BUS=="usb", ID=="9-1.2.7.3:1.0", SYMLINK+="gsmmd9"
BUS=="usb", ID=="9-1.2.7.4:1.0", SYMLINK+="gsmmd10"

#Mapping for USB Power Adaptor 2

BUS=="usb", ID=="9-1.1.1:1.0", SYMLINK+="gsmmd11"
BUS=="usb", ID=="9-1.1.2:1.0", SYMLINK+="gsmmd12"
BUS=="usb", ID=="9-1.1.3:1.0", SYMLINK+="gsmmd13"
BUS=="usb", ID=="9-1.1.4:1.0", SYMLINK+="gsmmd14"
BUS=="usb", ID=="9-1.1.5:1.0", SYMLINK+="gsmmd15"
BUS=="usb", ID=="9-1.1.6:1.0", SYMLINK+="gsmmd16"
BUS=="usb", ID=="9-1.1.7.1:1.0", SYMLINK+="gsmmd17"
BUS=="usb", ID=="9-1.1.7.2:1.0", SYMLINK+="gsmmd18"
BUS=="usb", ID=="9-1.1.7.3:1.0", SYMLINK+="gsmmd19"
BUS=="usb", ID=="9-1.1.7.4:1.0", SYMLINK+="gsmmd20"

#####Mapping for USB Power Adaptor 3#####

#3 USB ports which located top
BUS=="usb", ID=="2-1.5.1:1.0", SYMLINK+="gsmmd21"
BUS=="usb", ID=="2-1.5.2:1.0", SYMLINK+="gsmmd22"
BUS=="usb", ID=="2-1.5.3:1.0", SYMLINK+="gsmmd23"
BUS=="usb", ID=="2-1.5.4:1.0", SYMLINK+="gsmmd24"
BUS=="usb", ID=="2-1.5.5:1.0", SYMLINK+="gsmmd25"
BUS=="usb", ID=="2-1.5.6:1.0", SYMLINK+="gsmmd26"
BUS=="usb", ID=="2-1.5.7.1:1.0", SYMLINK+="gsmmd27"
BUS=="usb", ID=="2-1.5.7.2:1.0", SYMLINK+="gsmmd28"
BUS=="usb", ID=="2-1.5.7.3:1.0", SYMLINK+="gsmmd29"
BUS=="usb", ID=="2-1.5.7.4:1.0", SYMLINK+="gsmmd30"

BUS=="usb", ID=="2-1.6.1:1.0", SYMLINK+="gsmmd31"
BUS=="usb", ID=="2-1.6.2:1.0", SYMLINK+="gsmmd32"
BUS=="usb", ID=="2-1.6.3:1.0", SYMLINK+="gsmmd33"
BUS=="usb", ID=="2-1.6.4:1.0", SYMLINK+="gsmmd34"
BUS=="usb", ID=="2-1.6.5:1.0", SYMLINK+="gsmmd35"
BUS=="usb", ID=="2-1.6.6:1.0", SYMLINK+="gsmmd36"
BUS=="usb", ID=="2-1.6.7.1:1.0", SYMLINK+="gsmmd37"
BUS=="usb", ID=="2-1.6.7.2:1.0", SYMLINK+="gsmmd38"
BUS=="usb", ID=="2-1.6.7.3:1.0", SYMLINK+="gsmmd39"
BUS=="usb", ID=="2-1.6.7.4:1.0", SYMLINK+="gsmmd40"

##Mapping CDMA modem
#BUS=="usb", ID=="5-3:1.0", SYMLINK+="cdmamd"

##Mapping GSM71 Ex CDMA
BUS=="usb", ID=="5-3:1.0", SYMLINK+="gsmmd61"

##Mapping USB Modem via CARD
#BUS=="usb", ID=="6-2:1.0", SYMLINK+="gsmmd41"
BUS=="usb", ID=="6-2:1.0", SYMLINK+="gsmmd41"
BUS=="usb", ID=="5-2:1.0", SYMLINK+="gsmmd42"
BUS=="usb", ID=="6-1:1.0", SYMLINK+="gsmmd43"
BUS=="usb", ID=="8-1:1.0", SYMLINK+="gsmmd44"

#####Mapping for USB HUB ORION 10 Ports#####
BUS=="usb", ID=="2-1.1.4.2:1.0", SYMLINK+="gsmmd52"
BUS=="usb", ID=="2-1.1.4.3:1.0", SYMLINK+="gsmmd53"
BUS=="usb", ID=="2-1.1.4.4:1.0", SYMLINK+="gsmmd54"
#BUS=="usb", ID=="2-1.1.2:1.0", SYMLINK+="gsmmd55"
#BUS=="usb", ID=="2-1.1.1:1.0", SYMLINK+="gsmmd56"
BUS=="usb", ID=="2-1.1.4.2:1.0", SYMLINK+="gsmmd56"
BUS=="usb", ID=="2-1.1.3.2:1.0", SYMLINK+="gsmmd57"
BUS=="usb", ID=="2-1.1.3.1:1.0", SYMLINK+="gsmmd58"


#####Mapping for USB Power Adaptor 4####

#TESTING 
BUS=="usb", ID=="3-5.1:1.0", SYMLINK+="usb1"
BUS=="usb", ID=="3-5.2:1.0", SYMLINK+="usb2"
BUS=="usb", ID=="3-5.3:1.0", SYMLINK+="usb3"

